# adminstatus/crontab
* * * * * root /audit_required.sh
*/5 * * * * root /auto_delete.sh
