# adminstatus/crontab
* * * * * root /audit_required.sh
0 12 * * * root /auto_delete.sh
